# jetson/aeromed_tau_autopilot.service
[Unit]
Description=aeromed_tau_autopilot

[Service]
ExecStart=source /home/jetson/amt/jetson/autostart_service.sh

[Install]
WantedBy=default.target